# HOME/.config/systemd/user/notes.timer
[Unit]
Description=Sync notes repo periodically

[Timer]
OnStartupSec=10sec
OnUnitActiveSec=30min

[Install]
WantedBy=timers.target
