# HOME/.config/systemd/user/weather.timer
[Unit]
Description=Check the weather periodically

[Timer]
OnStartupSec=15sec
OnUnitActiveSec=15min

[Install]
WantedBy=timers.target
